# apps/data-collection/league-data-collection.service
[Unit]
Description=League Data Collection Service
After=network.target

[Service]
Type=forking
User=azureuser
WorkingDirectory=/home/azureuser/draftking-monorepo/apps/data-collection
Environment=HOME=/home/azureuser
Environment=NVM_DIR=/home/azureuser/.nvm
Environment=PATH=/home/azureuser/.nvm/versions/node/v22.11.0/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production

# Add debug flag to PM2
ExecStart=/home/azureuser/.nvm/versions/node/v22.11.0/bin/pm2 start ecosystem.config.cjs --no-daemon
ExecReload=/home/azureuser/.nvm/versions/node/v22.11.0/bin/pm2 reload all
ExecStop=/home/azureuser/.nvm/versions/node/v22.11.0/bin/pm2 stop all

# Modify restart behavior
Restart=on-failure
RestartSec=10s
StartLimitInterval=500s
StartLimitBurst=5

# Verbose logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target